# firstboot.service
[Unit]
Description=One time kiosk config boot script
[Service]
Type=simple
ExecStart=/opt/kiosk/initkio.sh
[Install]
WantedBy=multi-user.target 
